# BruteforceAnomaly/config/bruteforce-anomaly.service
[Unit]
Description=TrustLab SSH Brute Force Detection & APT Correlation System
After=network.target

[Service]
Type=simple
User=root
# These paths will be replaced with absolute paths by setup.sh
WorkingDirectory=/path/to/installation
ExecStart=/usr/bin/python3 /path/to/installation/main.py --monitor
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
