# BruteforceAnomaly/config/bruteforce-anomaly.service
[Unit]
Description=TrustLab SSH Brute Force Detection & APT Correlation System
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=%INSTALLATION_PATH%
ExecStart=%INSTALLATION_PATH%/config/trustlab_service.sh
ExecStop=%INSTALLATION_PATH%/config/stop_trustlab_service.sh
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
